# tools/transcribe_api/bandolier-transcribe_api-server.service

[Unit]
Description=bandolier-transcribe_api-server

[Service]
ExecStart=/home/sam/Documents/GitHub/Bandolier/tools/transcribe_api/venv-transcribe-api/bin/python /home/sam/Documents/GitHub/Bandolier/tools/transcribe_api/server.py
Restart=always
User=sam
Group=sam
WorkingDirectory=/home/sam/Documents/GitHub/Bandolier/tools/transcribe_api

[Install]
WantedBy=multi-user.target
